# templates/garage.service.j2
[Unit]
Description=Garage Data Store
Documentation=https://garagehq.deuxfleurs.fr/documentation/
After=network.target network-online.target systemd-user-sessions.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{ garage_home }}
{% if garage_env_variables | default({}, true) | length > 0 %}
EnvironmentFile=-{{ garage_config_dir }}/garage.env
{% endif %}
{% set _data_dirs = (garage_data_dirs | default([]) | length > 0) | ternary(garage_data_dirs, [garage_data_dir]) %}
{% set ns = namespace(paths=[]) %}
{% for item in _data_dirs %}
{% set ns.paths = ns.paths + [item.path if item is mapping else item] %}
{% endfor %}
ExecStart={{ garage_binary_path }} --config {{ garage_config_dir }}/garage.toml server
User={{ garage_user }}
Group={{ garage_group }}
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ garage_home }} {{ ns.paths | join(' ') }} {{ garage_metadata_dir }}
LimitNOFILE=42000

[Install]
WantedBy=multi-user.target
